# config/includes.chroot/etc/cron.daily/deb-rkhunter
#! /bin/sh
/usr/local/bin/rkhunter --pkgmgr dpkg -c --sk
