# config/includes.chroot/etc/cron.daily/run-rkhunter
#! /bin/sh
/usr/local/bin/rkhunter -c --sk
